# sverre-sync-hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>nl.sverre.sync-hourly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>sync-from-icloud.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/erikplantenga/Desktop/Sverre_Vakantie</string>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/tmp/sverre-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sverre-sync.err.log</string>
</dict>
</plist>
